# vpn/connman-vpn.service.in
[Unit]
Description=ConnMan VPN service

[Service]
Type=dbus
BusName=net.connman.vpn
ExecStart=@sbindir@/connman-vpnd -n
StandardOutput=null

[Install]
WantedBy=multi-user.target
