# src/connman-wait-online.service.in
[Unit]
Description=Wait for network to be configured by ConnMan
Requisite=connman.service
After=connman.service
Before=network-online.target
DefaultDependencies=no
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart=@sbindir@/connmand-wait-online
RemainAfterExit=yes

[Install]
WantedBy=network-online.target
