# src/connman.service.in
[Unit]
Description=Connection service
DefaultDependencies=false
Conflicts=shutdown.target
RequiresMountsFor=@localstatedir@/lib/connman
After=dbus.service network-pre.target
Before=network.target shutdown.target remote-fs-pre.target
Wants=network.target remote-fs-pre.target

[Service]
Type=dbus
BusName=net.connman
Restart=on-failure
ExecStart=@sbindir@/connmand -n
StandardOutput=null
ProtectHome=true

[Install]
WantedBy=multi-user.target
