# src/connman.service.in
[Unit]
Description=Connection service
Requires=dbus.socket
After=dbus.socket
Before=remote-fs.target

[Service]
Type=dbus
BusName=net.connman
Restart=on-failure
ExecStart=@prefix@/sbin/connmand -n
StandardOutput=null

[Install]
WantedBy=multi-user.target
